# satellite/systemd/voice-display-control.service
[Unit]
Description=Voice Assistant — Display MQTT Control Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dd
ExecStart=/usr/bin/python3 /home/dd/dev/voice-assistant/src/display/control.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
